# deploy_tools/firmware_finder.service
[Unit]
Description=Service For Firmware.TurkishJohnnys
After=network.target

[Service]
Restart=on-failure
User=jacob
WorkingDirectory=/home/jacob/sites/djacob/source/mysite
ExecStart=/home/jacob/sites/djacob/source/mysite /home/jacob/sites/djacob/.env/bin/gunicorn --bind unix:/tmp/firmware.turkishjohnnys.com.socket mysite.wsgi:application


[Install]
WantedBy=multi-user.target